# root/etc/s6-overlay/s6-rc.d/p11-libreoffice-customize/run
#!/usr/bin/with-contenv bash
set -e

# bring your own custom setup script
echo "Nothing to customize for LibreOffice.";
